# config/websockets.conf
[program:websockets]
command=php /var/www/vhosts/implementation.store/httpdocs/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=laravel-echo